# de.djermn.now-playing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.djermn.now-playing</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bariser/Library/Application Support/djermn/now_playing_daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bariser/Library/Logs/djermn-now-playing.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bariser/Library/Logs/djermn-now-playing.log</string>
</dict>
</plist>
